# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:dbus]
command=/usr/bin/dbus-daemon --system --nofork
autostart=true
autorestart=true
priority=5

[program:xvfb]
command=/usr/bin/Xvfb :0 -screen 0 %(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)sx24 -ac
autostart=true
autorestart=true
priority=10

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":0"
autostart=true
autorestart=true
priority=20

[program:x11vnc]
command=/usr/bin/x11vnc -display :0 -rfbauth /root/.vnc/passwd -forever -shared -rfbport 5900 -noxdamage
autostart=true
autorestart=true
priority=30

[program:novnc]
command=/usr/bin/novnc_proxy --vnc localhost:5900 --listen 6080
autostart=true
autorestart=true
priority=40

[program:antigravity]
command=/opt/scripts/start-antigravity.sh
environment=DISPLAY=":0"
autostart=true
autorestart=true
priority=50
startsecs=5
startretries=3
